# services.d/windshaft-cartodb/run
#!/usr/bin/with-contenv bash
exec 2>&1
cd /usr/src/app
exec node ./app.js
